# deploy/gunicorn_crm_metro.socket
[Unit]
Description=gunicorn socket for crm_metro

[Socket]
ListenStream=/run/gunicorn_crm_metro.sock

[Install]
WantedBy=sockets.target